# template/supervisord.conf
[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port=8101


[supervisord]
logfile={{DEPLOY_WORK_DIR}}/logs/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile={{DEPLOY_WORK_DIR}}/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
# user=investx                                        ; default user
childlogdir={{DEPLOY_WORK_DIR}}/logs               ; where child log files will live

[supervisorctl]
serverurl=http://0.0.0.0:8101


[program:{{APP_NAME}}]
command={{DEPLOY_WORK_DIR}}/venv/bin/gunicorn --bind 0.0.0.0:8100 project.wsgi --chdir={{APP_ROOT}} --env DJANGO_SETTINGS_MODULE=project.settings

directory={{DEPLOY_WORK_DIR}} ; all relative to this directory
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=TERM

startsecs=10

stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=12

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 60

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998